# deploy/devflow.service
[Unit]
Description=DevFlow FastAPI Backend
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/DevFlow/backend
EnvironmentFile=/home/ubuntu/DevFlow/.env
ExecStart=/home/ubuntu/DevFlow/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
